# systemd/discord-voice-niko-bridge.service
[Unit]
Description=Discord Voice Niko - Local Bridge (respond endpoint)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/.openclaw/workspace/discord-voice-niko
EnvironmentFile=/root/.openclaw/workspace/discord-voice-niko/.env
ExecStart=/usr/bin/node /root/.openclaw/workspace/discord-voice-niko/bridge/server.js
Restart=always
RestartSec=2

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
